# s6/docker-gen/run
#!/bin/sh

exec docker-gen -watch -only-exposed \
     -notify "s6-svc -h /s6/nginx" \
     /etc/nginx/nginx.tmpl /etc/nginx/conf.d/default.conf
